# docker/com.rxindexer.cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rxindexer.cleanup</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>cd /Users/rxindexer/Desktop/RXinDexer_1/docker && docker compose --profile maintenance run --rm cleanup</string>
    </array>
    
    <!-- Run every Sunday at 3:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/tmp/rxindexer-cleanup.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/rxindexer-cleanup.log</string>
    
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
